# deployer/deb/aggregator/control/postrm
#!/bin/sh
if [ "$1" = purge ]; then
	update-rc.d aggregator remove
	rm -rf /var/log/ourgrid/aggregator*.log*
fi

exit 0
